# build/templates/supervisor_marvin.conf
[program:marvin]
chdir=%(project_root)s
command=%(virtualenv_path)s/bin/marvin
master=True
vaccuum=True
max-requests=5000
autostart=true
autorestart=true
redirect_stderr=true
environment=LANG="en_US.UTF-8",LC_ALL="en_US.UTF-8",LC_LANG="en_US.UTF-8"
user=marvin
stopsignal=INT
